# confs/deploy.conf
start on runlevel [2345]
stop on runlevel [!2345]

respawn
setuid ubuntu
setgid ubuntu

script
    . /etc/environment
    /bin/bash /home/ubuntu/apps/{{ PROJECT_FILENAME }}/repository/run_on_server.sh fab $DEPLOYMENT_TARGET daemons.deploy >> {{ SERVER_LOG_PATH }}/deploy.log 2>&1
end script

post-stop exec sleep 30
